# deploy/supervisord.conf
[program:TaxApp]
command=/home/django/TaxApp/deploy/gunicorn.sh
directory=/home/django/TaxApp
user=nobody
autostart=true
autorestart=true
redirect_stderr=true